# backend/deploy/app.service
[Unit]
Description=Gunicorn instance to serve Crypto Volume Screener
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/var/www/crypto-screener
Environment="PATH=/var/www/crypto-screener/venv/bin"
ExecStart=/var/www/crypto-screener/venv/bin/gunicorn -w 4 -k uvicorn.workers.UvicornWorker main:app --bind 127.0.0.1:8000 --timeout 120

[Install]
WantedBy=multi-user.target
